# update-motd.d/00-figlet
#!/bin/sh

. /etc/update-motd.d/colors

printf "\n"$LIGHT_RED
#figlet "  "$(hostname -s)
toilet -f standard -F metal $(hostname -s)
printf $NONE
printf "\n"
